# aap/collections/ansible_collections/ansible/containerized_installer/roles/automationcontroller/templates/metrics-utility-build-report.timer.j2
[Unit]
Description=Run metrics-utility build report
Requires=metrics-utility-build-report.service

[Timer]
OnCalendar={{ metrics_utility_cronjob_report_schedule }}
Persistent=true

[Install]
WantedBy=timers.target